# scripts/exonus-tech.service
[Unit]
After=network.service

[Service]
ExecStart=/home/exonus-tech/config/exonus-script.sh

[Install]
WantedBy=default.target
